# s6-overlay/s6-rc.d/laravel-octane/run
#!/command/execlineb -P
with-contenv
backtick -E { mkdir -p /tmp/octane }
export OCTANE_STATE_PATH /tmp/octane/state.json
cd /var/www/html
php
artisan
octane:start
--server=frankenphp
--host=0.0.0.0
--port=8001
--workers=auto
--max-requests=1000
--watch