# s6-overlay/s6-rc.d/laravel-octane/run
#!/command/execlineb -P
with-contenv
importas -D "0.0.0.0" OCTANE_HOST OCTANE_HOST
importas -D "8001" OCTANE_PORT OCTANE_PORT
importas -D "2019" OCTANE_ADMIN_PORT OCTANE_ADMIN_PORT
export OCTANE_STATE_PATH /tmp/octane/state.json
cd /var/www/html
php
artisan
octane:start
--server=frankenphp
--host=${OCTANE_HOST}
--port=${OCTANE_PORT}
--admin-port=${OCTANE_ADMIN_PORT}
--workers=auto
--max-requests=1000
--watch