# s6-overlay/s6-rc.d/laravel-octane/run
#!/command/execlineb -P
with-contenv
cd /var/www/html
php
artisan
octane:frankenphp
--host=0.0.0.0
--port=8001
--workers=auto
--task-workers=auto
--max-requests=1000
--watch